# debian/cron.d
# on-twitterbot (export new Opennet Blog entries to Opennet twitter)
*/10 *	* * *	on-twitterbot	/usr/lib/on-twitterbot/twitterbot.py rss >/dev/null
